# prosody/rootfs/etc/s6-overlay/s6-rc.d/20-prosody/run
!/command/with-contenv bash

touch /tmp/xxx
sleep 600
